# usr/lib/systemd/system/greenboot-set-rollback-trigger.service
[Unit]
Description=Greenboot Rollback trigger
DefaultDependencies=no
After=local-fs.target
Before=systemd-update-done.service greenboot-healthcheck.service
Before=sysinit.target shutdown.target
Conflicts=shutdown.target
ConditionNeedsUpdate=|/etc
ConditionNeedsUpdate=|/var
RequiresMountsFor=/boot

[Service]
Type=oneshot
Restart=no
ExecStart=/usr/libexec/greenboot/greenboot set-rollback-trigger
PrivateMounts=yes

[Install]
WantedBy=systemd-update-done.service greenboot-healthcheck.service